# tools/owldrive_service/systemd/owldrive.service
[Unit]
Description=owlDrive CAN web interface
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=@SERVICE_DIR@
Environment=OWLDRIVE_CAN_CHANNEL=can0
Environment=OWLDRIVE_CAN_BITRATE=1000000
Environment=OWLDRIVE_HOST_NODE_ID=62
Environment=OWLDRIVE_CAN_MSG_ID=300
ExecStart=@SERVICE_DIR@/.venv/bin/uvicorn owldrive_service.app:app --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
